# deploy/systemd/mimmy.service
[Unit]
Description=Mimmy trading loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/mimmy
EnvironmentFile=-/opt/mimmy/.env
# Secrets Manager 사용 시 .env 대신 아래 두 줄로 교체:
# Environment="MIMMY_USE_SECRETS_MANAGER=true"
# Environment="MIMMY_SECRETS_ID=mimmy/prod"
ExecStart=/home/ubuntu/.local/bin/uv run --frozen mimmy loop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
