# deploy/systemd/mimmy.service
[Unit]
Description=Mimmy trading loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/mimmy
Environment="MIMMY_USE_SECRETS_MANAGER=true"
Environment="MIMMY_SECRETS_ID=mimmy/prod"
Environment="AWS_REGION=ap-northeast-2"
ExecStart=/home/ubuntu/.local/bin/uv run --frozen mimmy loop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
